# supervisord.conf
[program:bookmarks]
user=nobody
autostart=true
startsecs=0
command=python bookmark.py
